# files/hev-fsh.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2019 OpenWrt.org

START=90
